# roles/schedule_ansible/files/regular-ansible.timer
[Unit]
Description=Scheduled Ansible

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
