# roles/schedule_ansible/files/regular-ansible.service
[Unit]
Description=Scheduled Ansible

[Service]
Type=oneshot
User=jonny
Restart=on-failure
RestartSec=5m
ExecStart=/usr/local/bin/regular-ansible.sh
